# KJ360/com.kj360.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.kj360.server</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd ~/path/to/vault/KJ360 && ./kj360.sh --no-browser</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/kj360.out.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/kj360.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
<!--
INSTALLATION:

1. Edit the path above to point to your actual vault location
2. Copy to LaunchAgents folder:
   cp com.kj360.plist ~/Library/LaunchAgents/

3. Load the agent:
   launchctl load ~/Library/LaunchAgents/com.kj360.plist

4. To stop:
   launchctl unload ~/Library/LaunchAgents/com.kj360.plist

5. Check logs:
   tail -f /tmp/kj360.out.log
   tail -f /tmp/kj360.err.log
-->
