# KJ360/com.kj360.dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.kj360.dashboard</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd ~/path/to/vault/KJ360 && ./scripts/generate-dashboard.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <array>
        <!-- Monday through Friday at 08:30 -->
        <dict>
            <key>Weekday</key>
            <integer>1</integer>
            <key>Hour</key>
            <integer>8</integer>
            <key>Minute</key>
            <integer>30</integer>
        </dict>
        <dict>
            <key>Weekday</key>
            <integer>2</integer>
            <key>Hour</key>
            <integer>8</integer>
            <key>Minute</key>
            <integer>30</integer>
        </dict>
        <dict>
            <key>Weekday</key>
            <integer>3</integer>
            <key>Hour</key>
            <integer>8</integer>
            <key>Minute</key>
            <integer>30</integer>
        </dict>
        <dict>
            <key>Weekday</key>
            <integer>4</integer>
            <key>Hour</key>
            <integer>8</integer>
            <key>Minute</key>
            <integer>30</integer>
        </dict>
        <dict>
            <key>Weekday</key>
            <integer>5</integer>
            <key>Hour</key>
            <integer>8</integer>
            <key>Minute</key>
            <integer>30</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/tmp/kj360-dashboard.out.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/kj360-dashboard.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
<!--
INSTALLATION:

1. Edit the path above to point to your actual vault location
2. Copy to LaunchAgents folder:
   cp com.kj360.dashboard.plist ~/Library/LaunchAgents/

3. Load the agent:
   launchctl load ~/Library/LaunchAgents/com.kj360.dashboard.plist

4. Test manually:
   ./scripts/generate-dashboard.sh

5. Check logs:
   tail -f /tmp/kj360-dashboard.out.log
-->
